# bing_supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]
port=9092
;username =shiyi
;password =123456

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
username=shiyi
password=123456

[supervisord]
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=info                                  ; info, debug, warn, trace
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/www/station/logs/

[program:bing]
process_name=ip_bing
directory=/var/www/station/
command=sh runbing.sh
stopsignal=KILL
killasgroup=true
stopasgroup=true
autostart=true                ; supervisord守护程序启动时自动启动
autorestart=true              ; supervisord守护程序重启时自动重启
;redirect_stderr=true          ; 将stderr重定向到stdout
;stdout_logfile = /var/www/station/bing_ip.log

[supervisorctl]
username=shiyi
password=123456




