# supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]
port=9001
username = shiyi
password =  123456

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
username=shiyi
password= 123456


[program:webscan-8001]
process_name=ip_bing
directory=/var/www/station/
command=sh start.sh 8001
stopsignal=KILL
killasgroup=true
stopasgroup=true
autostart=true                ; supervisord守护程序启动时自动启动
autorestart=true              ; supervisord守护程序重启时自动重启
redirect_stderr=true          ; 将stderr重定向到stdout
stdout_logfile = /var/www/station/supervisord_bing_ip.log

[supervisorctl]
username=shiyi
password=123456




